# deploy/dev-team-dispatch.service
# Singleton, long-lived unit for the dev-team HTTP dispatch service.
#
#   sudo cp deploy/dev-team-dispatch.service /etc/systemd/system/
#   sudo systemctl enable --now dev-team-dispatch.service
#
# Reads /etc/dev-team/dev-team.env, which must define:
#   - DEV_TEAM_DISPATCH_TOKEN  the bearer token callers present as
#                              `Authorization: Bearer <token>` (a credential —
#                              keep it out of version control; rotate by
#                              editing the env file and restarting the unit)
#   - credentials: CLAUDE_CODE_OAUTH_TOKEN (Claude subscription, from
#                  `claude setup-token`) or ANTHROPIC_API_KEY (Claude API);
#                  submitted jobs run real agents.
#   - optionally GITHUB_TOKEN/GH_TOKEN for cloning private repos.
#
# The service binds the tailnet IP ONLY (100.116.115.8 below) so it is never
# reachable from the public internet — it authenticates every route except
# /health, but it also holds Claude credentials and executes agent code, so
# keep both the sandbox hardening and the tailnet-only bind. Jobs run one at a
# time (single-flight): the box has one shared Claude subscription.
[Unit]
Description=dev-team HTTP dispatch service (authenticated remote job runner)
After=network-online.target tailscaled.service
Wants=network-online.target
Requires=tailscaled.service

[Service]
Type=simple
User=devteam
Group=devteam
WorkingDirectory=/opt/dev-team
EnvironmentFile=/etc/dev-team/dev-team.env
# --dashboard-workspace makes dispatched jobs visible on the standing
# `dev-team --dashboard` process (which watches /opt/dev-team/workspace): each
# job also journals its events there (as its own run) and mirrors its assess
# report under audit/<job-id>/. Each job still runs in its own isolated
# workspace under /opt/dev-team/jobs/<id>.
ExecStart=/opt/dev-team/.venv/bin/dev-team --dispatch --host 100.116.115.8 --port 8738 --dashboard-workspace /opt/dev-team/workspace
Restart=on-failure
RestartSec=5

# Hardening (mirrors deploy/dev-team@.service).
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/dev-team
ProtectControlGroups=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectKernelLogs=yes
ProtectHostname=yes
ProtectClock=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
RemoveIPC=yes
LockPersonality=yes
SystemCallArchitectures=native
# Stronger still — but this unit RUNS jobs, which may use the in-process
# container sandbox (--sandbox, rootless podman/docker); those engines need
# user namespaces and clone/mount syscalls. Do NOT enable the following here
# unless this host never uses --sandbox (and then also run the whole process in
# an outer container/VM — see DEPLOYMENT.md 5d):
#   CapabilityBoundingSet=
#   RestrictNamespaces=yes
#   SystemCallFilter=@system-service
#   RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

[Install]
WantedBy=multi-user.target
